# config/crontab
0 0 * * * /app/analytics.py > /dev/stdout 2>&1
0 * * * * flock -xn /tmp/unfinish_task.lock -c '/app/collect_unfinish_task.py > /dev/stdout 2>&1'
